# extensions/fts++/zeitgeist-fts.service.in
[Unit]
Description=Zeitgeist full-text search indexer
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.gnome.zeitgeist.SimpleIndexer
ExecStart=@pkglibexecdir@/zeitgeist-fts
Slice=session.slice
